# aconfmgr/files/var/lib/pacman/local/corefreq-server-1.96.5-1/files
%FILES%
usr/
usr/bin/
usr/bin/corefreqd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/corefreqd.service

