# aconfmgr/files/var/lib/pacman/local/pamac-cli-11.6.0-3/files
%FILES%
usr/
usr/bin/
usr/bin/pamac
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pamac-offline-upgrade.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/pamac
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/pamac.fish
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/pamac.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/pamac.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_pamac

